# ansible/roles/minecraft/files/services/minecraft/backup.timer
[Unit]
Description=Minecraft server backup
[Timer]
OnCalendar=02:00
Persistent=true
[Install]
WantedBy=timers.target